# crontab
* * * * * zsh -c 'autoload compinit; compinit -u'
