# crontab/completion updater
* * * * * zsh -ic 'compinit -u -d $ZDOTDIR/cache/.zcompdump'
