# crontab/completion updater
* * * * * zsh -ic 'compinit -u $ZDOTDIR/cache/.zcompdump'
